# deployment/docker/stretch/config/xvfb.supervisord.conf
[program:Xvfb]
command=bash -c "/usr/bin/Xvfb $DISPLAY -screen 0 $XVFB_WHD"
